# deploy/supervisor.conf
; Supervisor configuration for the Santo Domingo House Price API
;
; SETUP INSTRUCTIONS:
;   1. Copy this file to /etc/supervisor/conf.d/predict-api.conf
;   2. Create log directory: sudo mkdir -p /var/log/predict-api
;   3. Reload supervisor:
;        sudo supervisorctl reread
;        sudo supervisorctl update
;        sudo supervisorctl start predict-api
;   4. Verify: curl http://localhost:8000/health
;
; NOTES:
;   - Update 'command' path if your virtualenv is named differently.
;   - Update 'directory' if the project is cloned to a different path.
;   - stopasgroup + killasgroup are required to also kill Uvicorn worker
;     child processes on restart; without them, workers become orphans.

[program:predict-api]

; Full path to uvicorn inside the virtualenv (avoids PATH issues)
command=/home/ubuntu/Project1/venv/bin/uvicorn api:app --host 0.0.0.0 --port 8000 --workers 2

; Project root — must match the directory where api.py lives so that
; 'from chatbot.chat import ...' resolves correctly
directory=/home/ubuntu/Project1

user=ubuntu

; Start automatically when supervisord starts (survives reboots)
autostart=true

; Restart the process if it exits unexpectedly
autorestart=true

; Kill all worker child processes cleanly on stop/restart
stopasgroup=true
killasgroup=true

; Logs — check these if the API fails to start
stdout_logfile=/var/log/predict-api/stdout.log
stderr_logfile=/var/log/predict-api/stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
